# chkccard_kill.cron
59 23 * * * admin /usr/bin/ph/chkccard_kill_cron.s
